# deb/postinst
#!/bin/bash

#
# Copyright NetApp Inc, 2021 All rights reserved
#

echo " --> set directory ownership"
chown -R harvest:harvest /opt/harvest
mkdir -p /var/log/harvest
mkdir -p /var/run/harvest
chown -R harvest:harvest /var/log/harvest
chown -R harvest:harvest /var/run/harvest

echo "  --> copying service"
cp /opt/harvest/service/harvest.service /etc/systemd/system/
chmod 664 /etc/systemd/system/harvest.service
if [ ! -e /opt/harvest/harvest.yml ]; then
    cp /opt/harvest/harvest.example.yml /opt/harvest/harvest.yml
    echo "  --> use default config file [/opt/harvest/harvest.yml]"
fi
if [ ! -e /usr/bin/harvest ]; then
    ln -s /opt/harvest/bin/harvest /usr/bin/harvest && echo " --> link created"
fi
systemctl daemon-reload
systemctl start harvest.service
systemctl enable harvest.service
echo "  --> harvest service started through systemctl"